# deb_template/name-arch-version/DEBIAN/postinst
#!/bin/bash

set -e

systemctl daemon-reload
systemctl enable DEB_PACKAGE.service
systemctl start DEB_PACKAGE.service
