# deb_template/name-arch-version/DEBIAN/prerm
#!/bin/bash

set -e

systemctl stop DEB_PACKAGE.service
systemctl disable DEB_PACKAGE.service
